# 3. 做点编码以外的事/server_files/supervisor/video.conf
[program:video]
command=pipenv run gunicorn z_movie.wsgi -b 127.0.0.1:9999
directory=/server/z_movie
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
